# templates/llm-output-dragonfly-requirepass-empty-detector/examples/good/03_systemd_strong.service
[Unit]
Description=DragonflyDB

[Service]
EnvironmentFile=/etc/dragonfly/secret.conf
ExecStart=/usr/bin/dragonfly --port 6379 --requirepass=${DRAGONFLY_PASS_2024_LongRandomValue}
