# templates/llm-output-saltstack-master-open-mode-true-detector/samples/bad/salt_master_systemd_open_mode.service
[Unit]
Description=The Salt Master Server

[Service]
Type=notify
ExecStart=/usr/bin/salt-master --open-mode

[Install]
WantedBy=multi-user.target
