# templates/llm-output-chronograf-no-auth-detector/examples/good/02_loopback.service
[Unit]
Description=Chronograf (loopback only -- single-operator dev box)

[Service]
ExecStart=/usr/bin/chronograf --host 127.0.0.1 --port 8888

[Install]
WantedBy=multi-user.target
