# templates/llm-output-mysql-skip-grant-tables-detector/samples/bad/systemd_unit.service
[Service]
Type=notify
User=mysql
ExecStart=/usr/sbin/mariadbd --skip_grant_tables --user=mysql
Restart=on-failure
