# templates/llm-output-consul-acl-default-policy-allow-detector/examples/bad/03_cli_flag.service
[Unit]
Description=Consul Agent

[Service]
ExecStart=/usr/local/bin/consul agent -server -bootstrap-expect=3 -default-policy=allow -data-dir=/opt/consul
Restart=on-failure

[Install]
WantedBy=multi-user.target
