# templates/llm-output-vault-disable-mlock-true-detector/examples/bad/01_systemd.service
[Unit]
Description=HashiCorp Vault
After=network-online.target

[Service]
User=vault
Group=vault
Environment="VAULT_DISABLE_MLOCK=true"
ExecStart=/usr/bin/vault server -config=/etc/vault.d/vault.hcl
Restart=on-failure

[Install]
WantedBy=multi-user.target
