# templates/llm-output-kibana-elasticsearch-ssl-verificationmode-none-detector/examples/bad/01_systemd.service
[Unit]
Description=Kibana
After=network-online.target

[Service]
User=kibana
Group=kibana
Environment="ELASTICSEARCH_HOSTS=https://es.internal:9200"
Environment="ELASTICSEARCH_SSL_VERIFICATIONMODE=none"
ExecStart=/usr/share/kibana/bin/kibana
Restart=on-failure

[Install]
WantedBy=multi-user.target
