# templates/llm-output-traefik-api-insecure-true-detector/examples/bad/06_systemd_unit.service
[Unit]
Description=Traefik with insecure API (BAD)
After=network.target

[Service]
ExecStart=/usr/local/bin/traefik --providers.file.filename=/etc/traefik/dyn.yaml --api.insecure=true --entrypoints.web.address=:80
Restart=always

[Install]
WantedBy=multi-user.target
