# templates/llm-output-dragonfly-requirepass-empty-detector/examples/bad/03_systemd_empty_pass.service
[Unit]
Description=DragonflyDB
After=network.target

[Service]
ExecStart=/usr/bin/dragonfly --port 6379 --requirepass=
Restart=always

[Install]
WantedBy=multi-user.target
